# wslenv.service
[Unit]
Description=WSL environment service
After=sshd.service

[Service]
Type=simple
ExecStart=/opt/site-lisp/init/wslenvstart.sh

[Install]
